# tools/runner/com.pk.x-runner.timer
[Unit]
Description=Declencheur periodique du PK X Runner

[Timer]
OnBootSec=2min
OnUnitActiveSec=15min
AccuracySec=30s
Persistent=true

[Install]
WantedBy=timers.target
